# deploy/node-monitor-api.service
[Unit]
Description=Node Monitor API (Bitcoin dashboard backend)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/projects/node-monitor
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/pi/projects/node-monitor/venv/bin/python backend/api_server.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
